# deploy/kasra.service
# Kasra — systemd service unit
# Place in /etc/systemd/system/kasra.service
# systemctl daemon-reload && systemctl enable kasra && systemctl start kasra

[Unit]
Description=Kasra AI Development Security Governance Platform
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=kasra
Group=kasra
WorkingDirectory=/opt/kasra

# Environment
Environment=KASRA_APP_DATABASE_URL=sqlite:////opt/kasra/data/kasra.db
Environment=KASRA_APP_API_KEY=CHANGE_ME_IN_PRODUCTION
Environment=KASRA_APP_JWT_SECRET=CHANGE_ME_IN_PRODUCTION
Environment=KASRA_APP_LOG_LEVEL=info
Environment=KASRA_AUDIT__ENABLED=true
Environment=KASRA_AUDIT__JSONL_PATH=/opt/kasra/data/kasra-audit.jsonl

# ExecStart
ExecStart=/opt/kasra/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8080 --workers 4
ExecReload=/bin/kill -HUP $MAINPID

# Restart
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Security
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/kasra/data

# Limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
